# deploy/zt-update-agent.service
[Unit]
Description=ZeroTrust Agent/Worker auto-update (polling fallback for SSM)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/zt-worker/deploy/update-agent.sh

StandardOutput=journal
StandardError=journal
