# deploy/zt-update-agent.timer
[Unit]
Description=ZeroTrust Agent/Worker update polling timer (hourly fallback)

[Timer]
# 부팅 후 10분 + 매시간. SSM이 닿지 않을 때 안전망.
OnBootSec=10min
OnUnitActiveSec=1h
Unit=zt-update-agent.service

[Install]
WantedBy=timers.target
